# spamassassin/services/sa-update-init/up
#!/command/with-contenv sh
sa-update || true
