# openarc/services/openarc/run
#!/command/with-contenv sh
exec openarc -p inet:0.0.0.0:8892 -c /etc/openarc.conf -f
