# openarc/services/openarc/run
#!/command/with-contenv sh
exec openarc -c /etc/openarc.conf -f
