# postfix/services/postfix/run
#!/command/with-contenv sh

postfix set-permissions

exec /usr/sbin/postfix -c /etc/postfix start-fg
