# postfix/services/postfix/run
#!/command/with-contenv sh

postfix set-permissions

exec postfix -c /etc/postfix start-fg 2>/dev/null
