# postfix/services/postfix/run
#!/command/with-contenv sh
exec /usr/sbin/postfix -c /etc/postfix start-fg
