# postfix/services/saslauthd/run
#!/command/with-contenv sh

saslauthd -m /var/run/saslauthd -a ldap -d 2>/dev/null
