# postfix/services/saslauthd/run
#!/command/with-contenv sh

saslauthd -a ldap -d -m /var/run/saslauthd
