# spamassassin/services/spamd/run
#!/command/with-contenv sh
exec spamd -u mail -g mail --syslog=stderr -i 0.0.0.0 --allowed-ips 0.0.0.0/0
