# postfix/services/syslogd/run
#!/command/with-contenv sh
exec syslogd -n -O /dev/stdout
